# supervisord.conf
[supervisord]
nodaemon = true

[program:demoniacmadness]
command = python3 tor_app.py
